# deploy/systemd/clinicmonitoring-backend.service
[Unit]
Description=ClinicMonitoring API (Django + Socket.IO / uvicorn ASGI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/clinicmonitoring/backend
EnvironmentFile=-/etc/clinicmonitoring.env
Environment=BIND_HOST=127.0.0.1
Environment=PORT=8000
Environment=DJANGO_SETTINGS_MODULE=config.settings
ExecStart=/opt/clinicmonitoring/backend/.venv/bin/python run_dev.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
